# packages/dirvishbot/dirvishbot-0.0.7.tar.gz/dirvishbot-0.0.7/resources/systemd/dirvishbot.service
[Unit]
Description=Dirvish Telegram Notifcation Bot
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
PIDFile=/var/run/dirvishbot.pid
ExecStart=/usr/local/bin/dirvishbot
Restart=always

[Install]
WantedBy=multi-user.target
