# packages/FSTA/FSTA-0.1.3.tar.gz/FSTA-0.1.3/maison.service
[Unit]
Description= Une maison qui obeit à la voie
After=network-online.target

[Service]
ExecStart=/opt/FSTA/maison.sh
User=root
KillMode=process

[Install]
WantedBy=multi-user.target