# packages/netsu-plugin-networkmanager/netsu-plugin-networkmanager-0.1.1.tar.gz/netsu-plugin-networkmanager-0.1.1/conf/netsu-nm.service
[Unit]
Description=netsu NetworkManager dependency
Requires=NetworkManager.service
After=NetworkManager.service
Before=netsud.service

[Service]
ExecStart=/usr/bin/true
RemainAfterExit=yes

[Install]
RequiredBy=netsud.service
